# airootfs/etc/xdg/autostart/cnchi.desktop
[Desktop Entry]
Type=Application
Name=Reborn Installer
Exec=/usr/bin/cnchi-start.sh
Icon=/usr/share/applications/cnchi.png
Terminal=True
